# getdata.cron
30 6 * * * LC_ALL=C.UTF-8 /usr/local/bin/python /opt/weather/src/getdata.py daily    >> /opt/weather/log/get_daily.log  2>&1
30 5 1 * * LC_ALL=C.UTF-8 /usr/local/bin/python /opt/weather/src/getdata.py history  >> /opt/weather/log/get_history.log  2>&1
30 7 * * * LC_ALL=C.UTF-8 /usr/local/bin/python /opt/weather/src/getdata.py forecast >> /opt/weather/log/get_forecast.log  2>&1
# 30 6 * * * LC_ALL=C.UTF-8 bash /opt/weather/getdata.sh -t daily    >> /opt/weather/log/get_daily.log  2>&1
# 30 5 1 * * LC_ALL=C.UTF-8 bash /opt/weather/getdata.sh -t history  >> /opt/weather/log/get_history.log  2>&1
# 30 7 * * * LC_ALL=C.UTF-8 bash /opt/weather/getdata.sh -t forecast >> /opt/weather/log/get_forecast.log  2>&1
*/1 * * * * echo $(date) >> /opt/weather/log/test.log
