# getdata.cron
30 5 * * *  sh /opt/weather/getdata.sh -t daily
30 5 1 * *  sh /opt/weather/getdata.sh -t history
30 6 * * *  sh /opt/weather/getdata.sh -t forecast
